# huggingface/supervisord.conf
# ─────────────────────────────────────────────────────────────────────────────
# supervisord.conf — process manager for the HF Spaces single container
#
# Process start order (by priority value, lower = earlier):
#   10  litellm  — AI proxy, must be up before backend tries to register models
#   20  backend  — Node.js management API (has built-in retry for LiteLLM)
#   30  nginx    — reverse proxy, last to start so upstreams are ready
#
# All stdout/stderr are forwarded to /dev/stdout and /dev/stderr so HF Spaces
# can capture logs in its built-in log viewer.
#
# Bug fixes applied:
#   [Bug1] Removed `user=user` — already non-root via Dockerfile USER; adding
#          it again causes supervisord to crash with "Can't drop privilege".
#
#   [BugA] environment= values must be on a SINGLE LINE or with the value
#          starting on the SAME LINE as the key (not on the next line).
#
#          Python's ConfigParser (used by supervisord) parses:
#            environment=
#                KEY="val"
#          as value = '\nKEY="val"'  (leading newline).
#
#          Supervisord's env parser then fails with:
#            Error: Format string '\nKEY=...' for 'environment' is badly
#            formatted: float argument required, not dict
#          See: https://github.com/Supervisor/supervisor/issues/1123
#
#          With the leading \n, supervisord REFUSES TO START — the whole
#          container crashes at boot. Even if it didn't crash, LITELLM_BASE_URL
#          would not be set, so backend/src/litellm.js would fall back to
#          "http://litellm:4000" (Docker Compose hostname), which is
#          unresolvable inside a single container. All LiteLLM API calls fail.
#
#          FIX: Keep the entire value on one line after the = sign.
# ─────────────────────────────────────────────────────────────────────────────

[supervisord]
nodaemon=true                       ; stay in foreground so Docker sees it as PID 1
logfile=/dev/null                   ; supervisord's own log → suppress (child logs are enough)
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid        ; /var/run is root-owned; use /tmp
; NOTE: do NOT add `user=` here. The container already runs as uid=1000 via
;       Dockerfile USER. Adding `user=` when already non-root crashes supervisord.

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

# ── LiteLLM proxy ─────────────────────────────────────────────────────────────
[program:litellm]
command=litellm --config /app/litellm/config.yaml --port 4000 --num_workers 2
directory=/app
autostart=true
autorestart=true
startsecs=15                        ; give LiteLLM time to initialise before health checks
startretries=5
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; All env vars (LITELLM_MASTER_KEY etc.) are inherited from entrypoint.sh via
; exec → supervisord env inheritance. No environment= needed here.

# ── Backend (Node.js management API) ──────────────────────────────────────────
[program:backend]
command=node /app/backend/src/index.js
directory=/app/backend
autostart=true
autorestart=true
startsecs=5
startretries=10
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; FIX [BugA]: environment= value MUST start on the same line as the key.
; Starting it on the next line causes supervisord to parse a leading \n, which
; makes its env formatter raise "badly formatted" and refuse to start.
;
; These three vars are always the same inside this container and different from
; the docker-compose defaults, so we override them explicitly. All other vars
; (LITELLM_MASTER_KEY, GATEWAY_PUBLIC_URL, JWT_SECRET, DB_PATH, LOG_LEVEL)
; are inherited from entrypoint.sh exports via supervisord's environment.
environment=LITELLM_BASE_URL="http://127.0.0.1:4000",PORT="3001",NODE_ENV="production"

# ── nginx reverse proxy ────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -c /app/huggingface/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
startsecs=3
startretries=5
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0